# CRPC/etc/systemd/system/crpc-prepare.service
[Unit]
Description=CRPC prepare dirs
Before=crpc-sweep.service

[Service]
Type=oneshot
User=raffaello
ExecStart=/usr/bin/bash -lc 'mkdir -p /tmp/hackrf_sweeps_text /tmp/tiles /tmp/tiles_proc /tmp/tiles_done /tmp/crpc_logs /tmp/yolo_runs'

[Install]
WantedBy=crpc.target
