# CRPC/etc/systemd/system/crpc-cleanup.timer
[Unit]
Description=Run CRPC cleanup every 10 min

[Timer]
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=timers.target
