# CRPC/etc/systemd/system/crpc-cleanup.service
[Unit]
Description=CRPC cleanup (size + age)

[Service]
Type=oneshot
User=raffaello
ExecStart=/usr/local/bin/crpc_cleanup.sh
